# platdriver/mydriveraddload.rules
ACTION=="add" ,KERNEL=="sample", RUN+="/sbin/modprobe -b platdev"
