# platdriver/mydriveraddremove.rules
ACTION=="remove" ,KERNEL=="sample", RUN+="/sbin/modprobe -b platdev"
